# ansible/files/udev/61-onbattery.rules
# Rule for when switching to battery
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="0",RUN+="/home/user/.config/scripts/udev.sh battery"
