# ansible/files/udev/62-onpower.rules
# Rule for when switching to power supply
SUBSYSTEM=="power_supply",ENV{POWER_SUPPLY_ONLINE}=="1",RUN+="/usr/bin/system76-power profile performance"
